# supervisor.mariadb.conf
[program:mariadb]
command=/usr/sbin/mysqld
priority=999
username=www-data
stdout_logfile=/tmp/%(program_name)s.stdout
stderr_logfile=/tmp/%(program_name)s.stderr


